# directory-integrity-monitor/systemd/directory-integrity-monitor.service
[Unit]
Description=Directory Integrity Monitor Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/directory-integrity-monitor
Restart=always
User=root

[Install]
WantedBy=multi-user.target
